# root/etc/profile.d/savapage.sh
pathmunge /opt/savapage/server/bin/linux-x64/
